# traffic_red_light_violation_detection-main_claude_fix/deploy/traffic-web.service
[Unit]
Description=Traffic Dashboard Web (Da Nang) - giam sat giao thong
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=nvidia
WorkingDirectory=/home/nvidia/traffic_project
ExecStart=/usr/bin/python3 /home/nvidia/traffic_project/web/server.py --port 8090
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
